# recipes-geco/m53/files/m53.service
[Unit]
Description=M53

[Service]
ExecStart=/usr/bin/m53
Restart=always

[Install]
WantedBy=multi-user.target