# recipes-alustrategy/alucount/files/88-net.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="148f", ATTR{idProduct}=="5370", RUN+="/etc/init.d/networking restart"